# mac/com.codex.thunder-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codex.thunder-bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>/Users/jiangguiqi/.thunder-bridge/thunder-agent.py</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>StandardOutPath</key>
  <string>/Users/jiangguiqi/.thunder-bridge/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/jiangguiqi/.thunder-bridge/launchd.err.log</string>
  <key>WorkingDirectory</key>
  <string>/Users/jiangguiqi/.thunder-bridge</string>
</dict>
</plist>
